# server/scripts/mapcache/mapcache.timer
[Unit]
Description=mapcache background job

[Timer]
OnUnitActiveSec=24h
OnBootSec=24h

[Install]
WantedBy=timers.target
